# kobo/KoboRoot/etc/udev/rules.d/97-bookshift.rules
KERNEL=="eth*", ACTION=="add", RUN+="/usr/local/bookshift/launcher.sh"
KERNEL=="wlan*", ACTION=="add", RUN+="/usr/local/bookshift/launcher.sh"
